# deploy/confessions.service
# Файл systemd-сервісу.
# Розташування на сервері: /etc/systemd/system/confessions.service
#
# Встановлення:
#   sudo cp deploy/confessions.service /etc/systemd/system/confessions.service
#   sudo systemctl daemon-reload
#   sudo systemctl enable confessions
#   sudo systemctl start confessions
#
# ВАЖЛИВО: шляхи нижче (/opt/confessions_board) мають відповідати
# реальному розташуванню проєкту на EC2-інстансі.

[Unit]
Description=Confessions Board (Django + Gunicorn)
After=network.target postgresql.service

[Service]
Type=simple
User=ubuntu
Group=www-data
WorkingDirectory=/opt/confessions_board
EnvironmentFile=/opt/confessions_board/.env
ExecStart=/opt/confessions_board/venv/bin/gunicorn \
    config.wsgi:application \
    --bind 0.0.0.0:8000 \
    --workers 3 \
    --access-logfile /var/log/confessions/access.log \
    --error-logfile /var/log/confessions/error.log
Restart=always
RestartSec=5

[Install]
WantedBy=multi-user.target
